# debian/postinst
#!/bin/sh

# Add rules copy init script
if [ -x "/etc/init.d/sci-udev-netrules.sh" ]; then
       update-rc.d sci-udev-netrules.sh defaults >/dev/null
fi

